# build/systemd_units/azure_cloud_manager_event_catcher@.service
[Unit]
PartOf=azure_cloud_manager_event_catcher.target
[Install]
WantedBy=azure_cloud_manager_event_catcher.target
[Service]
WorkingDirectory=/var/www/miq/vmdb
Environment=BUNDLER_GROUPS=manageiq_default,ui_dependencies
ExecStart=/bin/bash -lc 'exec ruby lib/workers/bin/run_single_worker.rb ManageIQ::Providers::Azure::CloudManager::EventCatcher --heartbeat --guid=%i'
Restart=no
Type=notify
Slice=miq-azure_cloud_manager_event_catcher.slice
